# private_dot_config/systemd/user/udiskie.service
[Unit]
PartOf=udiskie.service
After=udiskie.service
Requisite=udiskie.service

[Service]
ExecStart=udiskie
Restart=on-failure
